# launchd/com.calendarbridge.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.calendarbridge.sync</string>

  <!-- Use the wrapper so we only run when Outlook is actually up -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/Users/bel/calendarBridge/run_sync_wrapper.sh</string>
  </array>

  <!-- Run approximately every 30 minutes -->
  <key>StartInterval</key>
  <integer>1800</integer>

  <!-- Also run once right after login/load -->
  <key>RunAtLoad</key>
  <true/>

  <!-- LaunchAgent logs (in addition to wrapper logs) -->
  <key>StandardOutPath</key>
  <string>/Users/bel/calendarBridge/logs/launchd.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bel/calendarBridge/logs/launchd.stderr.log</string>

  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>
  <key>AbandonProcessGroup</key>
  <true/>
</dict>
</plist>
